# fever-screen.conf
[program:feverscreen]
command=/code/feverscreen/feverscreen
autostart=true
autorestart=true